# etc/systemd/system/build-archlinux-armv7h-docker.timer
[Unit]
Description=build-archlinux-armv7-docker.timer

[Timer]
OnCalendar=daily

[Install]
WantedBy=multi-user.target
