# etc/systemd/system/daily.target
[Unit]
Description=daily.target
StopWhenUnneeded=yes
