# etc/systemd/system/daily.timer
[Unit]
Description=daily.timer

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
